# deploy/aurora-subsurface.service
[Unit]
Description=Aurora Strata Subsurface Daemon
After=network-online.target sound.target
Wants=network-online.target
ConditionPathExists=__AURORA_ROOT__/aurora_subsurface_daemon.py

[Service]
Type=simple
User=__AURORA_USER__
Group=__AURORA_USER__
WorkingDirectory=__AURORA_ROOT__
Environment=HOME=__AURORA_HOME__
Environment=AURORA_SKIP_DEP_INSTALL=1
Environment=AURORA_DAEMON_BOOT_GREETING=1
Environment=AURORA_ALLOW_X_CLIENTS=0
Environment=AURORA_ENABLE_ROOM_OPERATOR=0
Environment=AURORA_ENABLE_QUIET_WINDOW=0
Environment=CLAUDE_TTS_VOICE=en-GB-SoniaNeural
Environment=AURORA_TTS_ROUTE=simple_first
Environment=PYTHONUNBUFFERED=1
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/__AURORA_UID__/bus
Environment=XDG_RUNTIME_DIR=/run/user/__AURORA_UID__
ExecStart=__AURORA_ROOT__/scripts/run_subsurface_daemon.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
